# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:xvfb]
command=Xvfb :99 -screen 0 %(ENV_SCREEN_WIDTH)sx%(ENV_SCREEN_HEIGHT)sx%(ENV_SCREEN_DEPTH)s
autorestart=true
priority=10
stdout_logfile=/var/log/xvfb.log
stderr_logfile=/var/log/xvfb_err.log

[program:fluxbox]
command=fluxbox
environment=DISPLAY=":99"
autorestart=true
priority=20
stdout_logfile=/var/log/fluxbox.log
stderr_logfile=/var/log/fluxbox_err.log

[program:x11vnc]
command=/bin/bash -c "sleep 2 && x11vnc -display :99 -nopw -listen 0.0.0.0 -forever -shared -rfbport 5900"
autorestart=true
priority=30
stdout_logfile=/var/log/x11vnc.log
stderr_logfile=/var/log/x11vnc_err.log

[program:websockify]
command=websockify --web /usr/share/novnc 6080 localhost:5900
autorestart=true
priority=40
stdout_logfile=/var/log/websockify.log
stderr_logfile=/var/log/websockify_err.log

[program:chromium]
command=/bin/bash -c "sleep 3 && chromium --no-sandbox --disable-gpu --disable-dev-shm-usage --no-zygote --single-process --no-first-run --disable-extensions --window-size=1280,800 --display=:99 about:blank"
environment=DISPLAY=":99"
autorestart=false
priority=50
stdout_logfile=/var/log/chromium.log
stderr_logfile=/var/log/chromium_err.log
